# config/files/usr/etc/systemd/system/u-firstboot.service
[Unit]
After=network-online.target
ConditionPathExists=!/var/lib/u-firstbooted

[Service]
User=root
Type=oneshot
ExecStart=/usr/lib/u-firstboot/install.sh
ExecStartPost=/usr/bin/touch /var/lib/u-firstbooted
RemainAfterExit=yes
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
